# deploy/qdrant-rag.service
# systemd unit for qdrant_rag — Phase 8b
#
# OPTIONAL alternative to `make up`. Default operator workflow uses Make.
# Install:
#   sudo cp deploy/qdrant-rag.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now qdrant-rag
#
# EDIT before installing:
#   - User= / Group=  → match your deploy user (must be in docker group)
#   - WorkingDirectory= → absolute path to the project root on this host
#   - EnvironmentFile= → absolute path to .env (must exist; chmod 600 recommended)

[Unit]
Description=qdrant_rag (multi-tenant Qdrant vector storage)
Requires=docker.service
After=docker.service network-online.target

[Service]
Type=simple
User=bol7
Group=bol7
WorkingDirectory=/home/bol7/Documents/BOL7/Qdrant
EnvironmentFile=-/home/bol7/Documents/BOL7/Qdrant/.env
ExecStartPre=/usr/bin/docker compose -f docker-compose.yml down --remove-orphans
ExecStart=/usr/bin/docker compose -f docker-compose.yml up
ExecStop=/usr/bin/docker compose -f docker-compose.yml down
Restart=on-failure
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
